# systemd/network-health-watcher.service
[Unit]
Description=Check network health and reboot if offline too long

[Service]
Type=oneshot
ExecStart=/opt/local/bin/network-health-watcher.sh